# lowbatteryd-open-rc-init-service
#!/sbin/runscript

start() {
	/usr/local/bin/lowbatteryd
	return 0
}

stop() {
       /usr/local/bin/lowbatteryd stop
       return 0
}